# config/supervisord.8.4/supervisord.conf
[supervisord]
user=root
nodaemon=true

[unix_http_server]
file=/run/supervisord.sock  ; the path to the socket file

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock ; use a unix:// URL for a unix socket

[program:php-fpm]
command=/bin/bash -c "mkdir -p /var/run/php && php-fpm84 --nodaemonize --fpm-config /etc/php84/php-fpm.conf"
user=root
autostart=true
autorestart=true
environment=LASER="1"
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:nginx]
command=/bin/bash -c "mkdir -p /run/nginx && sed s/__SUFFIX__/${SUFFIX}/g /etc/nginx/http.d/multiProjectLaravel.conf.template > /etc/nginx/http.d/multiProjectLaravel.conf && nginx -g 'daemon off;'"
user=root
autostart=true
autorestart=true
environment=LASER="1"
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
